# package/base-files/files/etc/hotplug.d/ntp/00-savelatest
[ "$ACTION" = "stratum" ] || exit 0

touch /var/state/ntp_last_sync
